# trading-bot-universal.service
[Unit]
Description=Trading Bot Universal
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bots/universal
ExecStart=/usr/bin/python3 /home/ubuntu/bots/bot2/universaltradingbot/universalbot.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
# Agregar límites de recursos para prevenir cuelgues por memoria
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target</content>
<parameter name="filePath">c:\Users\ACER\Documents\trading\trading-bot-universal.service